# automation-agent/deploy/packaging/deb/prerm
#!/bin/bash
set -e

# Stop service before removal
systemctl stop automation-agent.service || true
systemctl disable automation-agent.service || true

exit 0
